# systemd/claude-bot.service
[Unit]
Description=Claude Telegram bot
After=network-online.target
Wants=network-online.target
# Disable systemd's default crash-rate limit: we self-restart via /restart
# and the hour-watchdog, and we'd rather keep trying than ever stay down.
StartLimitIntervalSec=0

[Service]
Type=simple
User=claude
Group=claude
WorkingDirectory=/opt/claude-bot
EnvironmentFile=/etc/claude-bot.env
ExecStart=/usr/bin/node --enable-source-maps dist/index.js
Restart=on-failure
RestartSec=5s

NoNewPrivileges=true
ProtectSystem=strict
# ProtectHome=true shadows /home with an inaccessible tmpfs, and
# ReadWritePaths= does not punch through it (it only toggles writability on
# already-visible paths). The child CLI's cwd lives under /home/claude, so
# posix_spawn's chdir() fails → EACCES. Use tmpfs + BindPaths to bind the
# real /home/claude into the service's mount namespace.
ProtectHome=tmpfs
BindPaths=/home/claude
ReadWritePaths=/var/lib/claude-bot
PrivateTmp=true
RestrictSUIDSGID=true
LockPersonality=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
# RestrictNamespaces=true  # disabled: blocks the Claude CLI from spawning its subprocess (EACCES)
RestrictRealtime=true
MemoryMax=2G

[Install]
WantedBy=multi-user.target
